# docs/buildroot/rootfs_overlay/etc/systemd/system/pmc.service
[Unit]
Description=PMC (Payload Management Controller) System
After=network.target

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/usr/local/bin

# 启动命令
ExecStart=/usr/local/bin/pmc_control.sh start

# 停止命令
ExecStop=/usr/local/bin/pmc_control.sh stop

# 重启策略
Restart=on-failure
RestartSec=10

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pmc

[Install]
WantedBy=multi-user.target
